# systemd/alarm-builder@.timer
[Unit]
Description=Build ArchLinuxARM Daily for project %i

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target